# unified-chat.service.example
[Unit]
Description=Unified Chat
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<user>
Group=<user>
WorkingDirectory=/home/<user>/unified-chat
EnvironmentFile=/home/<user>/unified-chat/.env
ExecStart=/home/<user>/unified-chat/venv/bin/uvicorn unified_chat.main:app --host 0.0.0.0 --port 8090
Restart=always
RestartSec=5
TimeoutStopSec=20
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target

